# conf/supervisor.ini
[program:draftdotcloud]
command=/home/draftdotcloud/draft.cloud/start.sh
directory=/home/draftdotcloud/draft.cloud
user=draftdotcloud

